# nginx_supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
environment=
    http_proxy=
    https_proxy=
    HTTP_PROXY=
    HTTPS_PROXY=

